# install/backup.cron
00 * * * * /usr/bin/rlcraft/backup.sh -c -i /usr/bin/rlcraft/world -o /media/Backups -s rlcraft
