# install/whitelist.cron
*/5 * * * * /usr/bin/rlcraft/whitelist.sh
